# crontab
# SCP-JP Scripts Cron Jobs
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

# new_page/tagging: 10分ごと
*/10 * * * * cd /app && uv run scripts/new_page/tagging.py >> /var/log/scp-jp-scripts/new_page_tagging.log 2>&1

# collab_deletion/notice: 毎月1日 04:00 JST
0 4 1 * * cd /app && uv run scripts/collab_deletion/notice.py >> /var/log/scp-jp-scripts/collab_deletion_notice.log 2>&1

# collab_deletion/exec: 毎月4日 04:00 JST
0 4 4 * * cd /app && uv run scripts/collab_deletion/exec.py >> /var/log/scp-jp-scripts/collab_deletion_exec.log 2>&1
